# hb-ipsc-supervise.conf
[program:HB]
command = /opt/hblink/HB_Bridge.py


[program:IPSC]
command = /opt/dmrlink/bridge/IPSC_Bridge.py


